# files/base/etc/local.d/00sysctl.start
#!/bin/sh
#
# Temporary workaround since sysctl isn't loading sysctl.conf at boot

/sbin/sysctl --system
/sbin/sysctl -p /etc/sysctl.confl

