# files/mobile/etc/profile.d/mobile.custom.sh
###########################################
#                                         #
# MANAGED BY PUPPET                       #
#                                         #
# Manual changes WILL be overwritten      #
#                                         #
###########################################

ANDROID_HOME="/opt/android-sdk-update-manager"
export ANDROID_HOME

