# files/mobile/etc/profile.d/path.custom.sh
###########################################
#                                         #
# MANAGED BY PUPPET                       #
#                                         #
# Manual changes WILL be overwritten      #
#                                         #
###########################################

PATH="/opt/android-sdk-update-manager/platform-tools/:${PATH}" 
PATH="/opt/android-sdk-update-manager/tools/:${PATH}"

export PATH
