# system-infrastructure-puppet/templates/vcs/etc/xinetd.d/fossil.erb
###########################################
#                                         #
# MANAGED BY PUPPET                       #
#                                         #
# Manual changes WILL be overwritten      #
#                                         #
###########################################

service fossil
{
        type = UNLISTED
        socket_type = stream
        protocol = tcp
        wait = no
        user = deployer
        group = www-data
        disable = no
        server = /usr/bin/fossil
        # We only expect connections from localhost since external requests are being proxied by nginx to this fossil service
        only_from = 127.0.0.1 localhost

        # HTTP port to listen to (the same as the one in '/etc/services')
        port = 4545
        # The path that the Fossil repos are stored
        <% if @vcs_type == "server" %>
        server_args = http /srv/fossil/fossils --https --host fossil.internal.nitelite.io --notfound https://fossil.internal.nitelite.io
        <% end %>

        <% if @vcs_type == "remote" %>
        server_args = http /srv/fossil/fossils --https --host fossil.hackbytes.io --notfound https://fossil.hackbytes.io
        <% end %>

        <% if @vcs_type == "mirror" %>
          <% if @mirror_environment == "production" %>
            server_args = http /srv/fossil/fossils --https --host fossil.tehpotatoking.com --notfound https://fossil.tehpotatoking.com
          <% else %>
            server_args = http /srv/fossil/fossils --https --host <%= @mirror_environment %>.fossil.tehpotatoking.com --notfound https://<%= @mirror_environment %>.fossil.tehpotatoking.com
          <% end %>
        <% end %>
}
