# files/xorg-server/etc/local.d/10redshift.start
#!/bin/sh

###########################################
#                                         #
# MANAGED BY PUPPET                       #
#                                         #
# Manual changes WILL be overwritten      #
#                                         #
###########################################

if [ -f "/etc/redshift.conf" ];
then
  redshift -c /etc/redshift.conf &
fi
