# system-infrastructure-puppet/files/data/etc/cron.weekly/elog-cleanup
#!/bin/sh

find /var/log/portage/elog -type f -mtime +30 -exec rm '{}' \;
