# files/base/etc/profile.d/ids.custom.sh
###########################################
#                                         #
# MANAGED BY PUPPET                       #
#                                         #
# Manual changes WILL be overwritten      #
#                                         #
###########################################

# IDS for any logins
pinky | mail -s "[nitelite.io] Login on ${HOSTNAME}" bfs5030+nitelite@gmail.com

